# package/base-files/iopsys/etc/hotplug.d/iface/01-fap
[ ifup = "$ACTION" ] && {
	fapctl disable
	fcctl disable
	fcctl flush
	fcctl enable
	fapctl enable
}

