# package/base-files/iopsys/etc/hotplug.d/switch/05-fap-tm
case "$ACTION" in
	add|register)
		case $INTERFACE in
			eth*)
				fap tm --enable --if $INTERFACE --auto
				fap tm --apply --if $INTERFACE
			;;
		esac
	;;
	remove|unregister)
	;;
esac

