# package/base-files/files/etc/hotplug.d/net/01-bridge
include /lib/network

[ ${INTERFACE:0:2} == "wl" ] && exit

case "$ACTION" in
        add|register)
		case "$INTERFACE" in
			# disable ipv6 if it is an unmanaged bridge
			br-*)
				local net=${INTERFACE:3}
				local proto=$(uci -q get network.$net.proto)
				if [ -z "$proto" -o "$proto" == "none" ]; then
					echo 1 > /proc/sys/net/ipv6/conf/$INTERFACE/disable_ipv6
				fi
			;;
		esac
        ;;
esac

