# package/base-files/files/etc/hotplug.d/net/01-macfix
# if an uplink device gets the mac address of the basemac
# change it to use the mac address of its originating interface 

case "$ACTION" in
	add|register)
		case "$INTERFACE" in
			 [aep]t[hm][0-9]*.[0-9]*)
				local basemac ifmac baseif bifmac
				basemac=$(ifconfig bcmsw | grep HWaddr | awk '{print$NF}')
				ifmac=$(ifconfig $INTERFACE | grep HWaddr | awk '{print$NF}')
				if [ "$ifmac" == "$basemac" ]; then
					baseif=${INTERFACE:0:4}
					bifmac=$(ifconfig $baseif | grep HWaddr | awk '{print$NF}')
					ifconfig $INTERFACE hw ether $bifmac 2>/dev/null
				fi
			;;
		esac
	;;
esac

