# configs/eg300/fs/etc/init.d/eg300_wan_monitor
#!/bin/sh /etc/rc.common

START=99

boot()
{
    if [ "$(db get hw.board.hardware)" == "EG300" ]
    then
        /sbin/monitorwan&
    fi
}
