# package/zeta/basefiles/iopsys/base-files-iopsys/etc/hotplug.d/net/01-bridge
include /lib/network

[ ${INTERFACE:0:2} == "wl" ] && exit

case "$ACTION" in
        add|register)
		case "$INTERFACE" in
			br-*)
				local net=${INTERFACE:3}
				local proto=$(uci -q get network.$net.proto)
				# disable ipv6 if it is an unmanaged bridge
				if [ -z "$proto" -o "$proto" == "none" ]; then
					echo 1 > /proc/sys/net/ipv6/conf/$INTERFACE/disable_ipv6
				fi
			;;
		esac
        ;;
        remove|unregister)
		# reload a bridge if a wan interface gets removed from it
                if [ "$INTERFACE" != "${INTERFACE/./}" ]; then
                        networkname=$(get_network_of "$INTERFACE")
                        for net in $networkname; do
                                type=$(uci get network."$net".type)
                                if [ "$type" == "bridge" ]; then
                                        /sbin/ifup "$net"
                                fi
                        done
                fi
        ;;
esac

